# net-misc/lldpad/files/lldpad.init
#!/sbin/openrc-run

depend() {
    need net
}

start() {
  ebegin "Starting lldpad daemon"
  start-stop-daemon --start --pidfile /run/lldpad.pid --exec /usr/sbin/lldpad -- -d
  eend $?
}

stop() {
  ebegin "Stopping lldpad daemon"
  start-stop-daemon --stop --pidfile /run/lldpad.pid --exec /usr/sbin/lldpad --
  eend $?
}

